# diymount_0.0.2/DEBIAN/postinst
#!/bin/sh
systemctl enable diymount.service
